# scripts/launchd/com.mitchell.career-ops.telegram-bot.plist
<?xml version="1.0" encoding="UTF-8"?>
<!DOCTYPE plist PUBLIC "-//Apple//DTD PLIST 1.0//EN" "http://www.apple.com/DTDs/PropertyList-1.0.dtd">
<plist version="1.0">
<dict>
  <key>Label</key>
  <string>com.mitchell.career-ops.telegram-bot</string>

  <key>ProgramArguments</key>
  <array>
    <string>/Users/mitchellwilliams/.nvm/versions/node/v24.14.0/bin/node</string>
    <string>/Users/mitchellwilliams/Documents/career-ops/telegram-bot.mjs</string>
  </array>

  <key>WorkingDirectory</key>
  <string>/Users/mitchellwilliams/Documents/career-ops</string>

  <key>RunAtLoad</key>
  <true/>

  <key>KeepAlive</key>
  <true/>

  <key>StandardOutPath</key>
  <string>/Users/mitchellwilliams/Documents/career-ops/data/logs/telegram-bot.log</string>

  <key>StandardErrorPath</key>
  <string>/Users/mitchellwilliams/Documents/career-ops/data/logs/telegram-bot.log</string>

  <key>ThrottleInterval</key>
  <integer>10</integer>
</dict>
</plist>
