# scripts/launchd/com.mitchell.career-ops.weekly-light.plist
<?xml version="1.0" encoding="UTF-8"?>
<!DOCTYPE plist PUBLIC "-//Apple//DTD PLIST 1.0//EN" "http://www.apple.com/DTDs/PropertyList-1.0.dtd">
<plist version="1.0">
<dict>
    <key>Label</key>
    <string>com.mitchell.career-ops.weekly-light</string>

    <!-- Zero-LLM-cost weekly health run: verify-pipeline + analyze-patterns + followup-cadence -->
    <!-- Runs Saturday 08:10 PT — after batch quota reset, before human review window -->
    <!-- No API spend. Safe to run even on low-quota days. -->
    <key>ProgramArguments</key>
    <array>
        <string>/Users/mitchellwilliams/.nvm/versions/node/v24.14.0/bin/node</string>
        <string>/Users/mitchellwilliams/Documents/career-ops/scripts/weekly-light.mjs</string>
    </array>

    <key>WorkingDirectory</key>
    <string>/Users/mitchellwilliams/Documents/career-ops</string>

    <!-- Saturday 08:10 PT (Weekday 6 = Saturday) -->
    <key>StartCalendarInterval</key>
    <dict>
        <key>Weekday</key>
        <integer>6</integer>
        <key>Hour</key>
        <integer>8</integer>
        <key>Minute</key>
        <integer>10</integer>
    </dict>

    <key>RunAtLoad</key>
    <false/>

    <key>StandardOutPath</key>
    <string>/Users/mitchellwilliams/Documents/career-ops/data/logs/weekly-light.out</string>

    <key>StandardErrorPath</key>
    <string>/Users/mitchellwilliams/Documents/career-ops/data/logs/weekly-light.err</string>

    <key>EnvironmentVariables</key>
    <dict>
        <key>PATH</key>
        <string>/Users/mitchellwilliams/.nvm/versions/node/v24.14.0/bin:/usr/local/bin:/opt/homebrew/bin:/usr/bin:/bin</string>
        <key>HOME</key>
        <string>/Users/mitchellwilliams</string>
    </dict>
</dict>
</plist>
