# scripts/launchd/com.mitchell.career-ops.network-database-build.plist
<?xml version="1.0" encoding="UTF-8"?>
<!DOCTYPE plist PUBLIC "-//Apple//DTD PLIST 1.0//EN" "http://www.apple.com/DTDs/PropertyList-1.0.dtd">
<plist version="1.0">
<dict>
    <key>Label</key>
    <string>com.mitchell.career-ops.network-database-build</string>

    <!-- ZETA recurring: rebuild the canonical data/network-database.json from on-disk LinkedIn/X/contacts sources. -->
    <!-- Pure aggregation, $0/run. Daily 02:30 PT. Underlying agent: scripts/build-network-database.mjs (shipped by ZETA overnight 2026-05-19). -->
    <!-- Cadence guard: always (launchd handles the daily cadence). -->
    <key>ProgramArguments</key>
    <array>
        <string>/bin/bash</string>
        <string>/Users/mitchellwilliams/Documents/career-ops/scripts/wrappers/cron-run.sh</string>
        <string>network-database-build</string>
        <string>always</string>
        <string>/Users/mitchellwilliams/.nvm/versions/node/v24.14.0/bin/node</string>
        <string>scripts/build-network-database.mjs</string>
    </array>

    <key>WorkingDirectory</key>
    <string>/Users/mitchellwilliams/Documents/career-ops</string>

    <!-- Daily 02:30 PT -->
    <key>StartCalendarInterval</key>
    <dict>
        <key>Hour</key>
        <integer>2</integer>
        <key>Minute</key>
        <integer>30</integer>
    </dict>

    <key>RunAtLoad</key>
    <false/>

    <key>StandardOutPath</key>
    <string>/Users/mitchellwilliams/Documents/career-ops/data/logs/network-database-build.out</string>

    <key>StandardErrorPath</key>
    <string>/Users/mitchellwilliams/Documents/career-ops/data/logs/network-database-build.err</string>

    <key>EnvironmentVariables</key>
    <dict>
        <key>PATH</key>
        <string>/Users/mitchellwilliams/.nvm/versions/node/v24.14.0/bin:/usr/local/bin:/opt/homebrew/bin:/usr/bin:/bin</string>
        <key>HOME</key>
        <string>/Users/mitchellwilliams</string>
    </dict>
</dict>
</plist>
