# scripts/launchd/com.mitchell.career-ops.community-scan.plist
<?xml version="1.0" encoding="UTF-8"?>
<!DOCTYPE plist PUBLIC "-//Apple//DTD PLIST 1.0//EN" "http://www.apple.com/DTDs/PropertyList-1.0.dtd">
<plist version="1.0">
<dict>
  <key>Label</key>
  <string>com.mitchell.career-ops.community-scan</string>

  <key>ProgramArguments</key>
  <array>
    <string>/bin/zsh</string>
    <string>-lc</string>
    <string>exec "$NVM_DIR/nvm-exec" node /Users/mitchellwilliams/Documents/career-ops/community-scan.mjs</string>
  </array>

  <key>EnvironmentVariables</key>
  <dict>
    <key>NVM_DIR</key>
    <string>/Users/mitchellwilliams/.nvm</string>
  </dict>

  <key>WorkingDirectory</key>
  <string>/Users/mitchellwilliams/Documents/career-ops</string>

  <!-- Run at 02:30 PT daily, right after the portal scan at 02:00 -->
  <key>StartCalendarInterval</key>
  <dict>
    <key>Hour</key>
    <integer>10</integer>
    <key>Minute</key>
    <integer>30</integer>
  </dict>

  <key>StandardOutPath</key>
  <string>/Users/mitchellwilliams/Documents/career-ops/data/logs/community-scan.log</string>

  <key>StandardErrorPath</key>
  <string>/Users/mitchellwilliams/Documents/career-ops/data/logs/community-scan.log</string>
</dict>
</plist>
