# scripts/launchd/com.mitchell.career-ops.chrome-debugging.plist
<?xml version="1.0" encoding="UTF-8"?>
<!DOCTYPE plist PUBLIC "-//Apple//DTD PLIST 1.0//EN" "http://www.apple.com/DTDs/PropertyList-1.0.dtd">
<plist version="1.0">
<dict>
    <key>Label</key>
    <string>com.mitchell.career-ops.chrome-debugging</string>

    <!-- Keep a dedicated Chrome instance running on
         http://127.0.0.1:9222 with the profile at
         ~/Library/Application Support/career-ops-chrome-debug/.
         Playwright agents (phase-B-prime-mechanical-enrich.mjs,
         refresh-master per-contact handler) attach to this Chrome
         over CDP, inheriting the LinkedIn login Mitchell set up
         once via `node scripts/launch-debug-chrome.mjs --setup`. -->

    <key>ProgramArguments</key>
    <array>
        <string>/Users/mitchellwilliams/.nvm/versions/node/v24.14.0/bin/node</string>
        <string>/Users/mitchellwilliams/Documents/career-ops/scripts/launch-debug-chrome.mjs</string>
    </array>

    <key>WorkingDirectory</key>
    <string>/Users/mitchellwilliams/Documents/career-ops</string>

    <!-- Fire at login + keep alive. macOS Tahoe KeepAlive bug:
         after bootstrap, also run `launchctl start <label>` once
         to actually kick the process. RunAtLoad alone doesn't
         fire on Tahoe (per project memory note). -->
    <key>RunAtLoad</key>
    <true/>

    <key>KeepAlive</key>
    <true/>

    <key>ThrottleInterval</key>
    <integer>60</integer>

    <key>StandardOutPath</key>
    <string>/Users/mitchellwilliams/Documents/career-ops/data/logs/chrome-debugging.out</string>

    <key>StandardErrorPath</key>
    <string>/Users/mitchellwilliams/Documents/career-ops/data/logs/chrome-debugging.err</string>

    <key>EnvironmentVariables</key>
    <dict>
        <key>PATH</key>
        <string>/Users/mitchellwilliams/.nvm/versions/node/v24.14.0/bin:/usr/local/bin:/opt/homebrew/bin:/usr/bin:/bin</string>
        <key>HOME</key>
        <string>/Users/mitchellwilliams</string>
    </dict>
</dict>
</plist>
