# scripts/launchd/com.mitchell.career-ops.detector-health.plist
<?xml version="1.0" encoding="UTF-8"?>
<!DOCTYPE plist PUBLIC "-//Apple//DTD PLIST 1.0//EN" "http://www.apple.com/DTDs/PropertyList-1.0.dtd">
<plist version="1.0">
<dict>
    <key>Label</key>
    <string>com.mitchell.career-ops.detector-health</string>

    <key>ProgramArguments</key>
    <array>
        <string>/bin/zsh</string>
        <string>-lc</string>
        <string>exec "$NVM_DIR/nvm-exec" node /Users/mitchellwilliams/Documents/career-ops/scripts/agents/detector-health-check.mjs</string>
    </array>

    <key>EnvironmentVariables</key>
    <dict>
        <key>NVM_DIR</key>
        <string>/Users/mitchellwilliams/.nvm</string>
    </dict>

    <key>WorkingDirectory</key>
    <string>/Users/mitchellwilliams/Documents/career-ops</string>

    <!-- Sunday 08:00 PT — weekly detector health sweep per δ.NH.2 -->
    <key>StartCalendarInterval</key>
    <dict>
        <key>Weekday</key>
        <integer>0</integer>
        <key>Hour</key>
        <integer>8</integer>
        <key>Minute</key>
        <integer>0</integer>
    </dict>

    <key>RunAtLoad</key>
    <false/>

    <key>StandardOutPath</key>
    <string>/Users/mitchellwilliams/Documents/career-ops/data/logs/detector-health.out</string>

    <key>StandardErrorPath</key>
    <string>/Users/mitchellwilliams/Documents/career-ops/data/logs/detector-health.err</string>

    <key>EnvironmentVariables</key>
    <dict>
        <key>PATH</key>
        <string>/Users/mitchellwilliams/.nvm/versions/node/v24.14.0/bin:/usr/local/bin:/opt/homebrew/bin:/usr/bin:/bin</string>
        <key>HOME</key>
        <string>/Users/mitchellwilliams</string>
    </dict>
</dict>
</plist>
