# systemd/pipeline-orchestrator.service
[Unit]
Description=Yash Pipeline Orchestrator
After=network-online.target
Wants=network-online.target

[Service]
Type=simple
WorkingDirectory=/yash-superClaudeHuman/projects/yash-ai-automation-career
EnvironmentFile=/etc/yash-pipeline/agent.env
ExecStart=/usr/bin/env node services/pipeline-orchestrator.mjs
Restart=on-failure
RestartSec=5
WatchdogSec=300
MemoryMax=512M
StandardOutput=journal
StandardError=journal

[Install]
WantedBy=default.target
